# config/templates/mariadb-container.service.tmpl
# ===========================================================================
# slurm-mariadb.service — MariaDB container for Slurm accounting
# Generated by slurm-installer on %%GENERATED_DATE%%
# ===========================================================================
# This service manages the MariaDB container used by slurmdbd.
# Data is persisted in /var/lib/slurm-mariadb.
# ===========================================================================

[Unit]
Description=MariaDB Container for Slurm Accounting
After=network-online.target
Wants=network-online.target
Before=slurmdbd.service

[Service]
Type=simple
Restart=always
RestartSec=5
TimeoutStartSec=300

# Container runtime and name
Environment=CONTAINER_RUNTIME=%%CONTAINER_RUNTIME%%
Environment=CONTAINER_NAME=slurm-mariadb

# Start: run container (remove existing first to ensure clean state)
ExecStartPre=-/usr/bin/${CONTAINER_RUNTIME} rm -f ${CONTAINER_NAME}
ExecStart=/usr/bin/${CONTAINER_RUNTIME} run \
    --name ${CONTAINER_NAME} \
    --network=host \
    --volume /var/lib/slurm-mariadb:/var/lib/mysql:Z \
    --env MARIADB_ROOT_PASSWORD=%%MARIADB_ROOT_PASSWORD%% \
    --env MARIADB_DATABASE=%%MARIADB_DATABASE%% \
    --env MARIADB_USER=%%MARIADB_USER%% \
    --env MARIADB_PASSWORD=%%MARIADB_PASSWORD%% \
    docker.io/library/mariadb:lts

# Stop: stop the container
ExecStop=/usr/bin/${CONTAINER_RUNTIME} stop -t 30 ${CONTAINER_NAME}

[Install]
WantedBy=multi-user.target
